# deploy/systemd/homesite-backend.service
[Unit]
Description=HomeSite v2 Backend (FastAPI)
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=exec
User=homesite
Group=homesite
WorkingDirectory=/opt/homesite/backend
EnvironmentFile=/opt/homesite/.env
ExecStart=/opt/homesite/venv/bin/gunicorn app.main:app \
    --bind 127.0.0.1:8000 \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
